# .github/deb/postinst
#!/bin/bash
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable rosadisk-agent
    systemctl start rosadisk-agent
fi
